# otel-sync/crontab
* * * * * echo "Crontab is working - watchdog 1"
* * * * * node dist/index.js > /dev/stdout

